# docker-files/supervisor/crond.conf
[program:cron]
command = cron -f -L 15
startsecs = 2
stopwaitsecs = 5
stopasgroup = false
killasgroup = true
