# docker-files/supervisor/mysqld.conf
[program:mysqld]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
startsecs=2

