# docker-files/supervisor/rails.conf
[program:rails]
command=/root/rails.sh
autostart=true
autorestart=true
