# docker-files/supervisor/rails.conf
[program:rails]
command=sh -c "cd /srv/linuxfr && rails server thin"
autostart=true
autorestart=true
