# docker-files/supervisor/rsyslogd.conf
[program:rsyslog]
command = rsyslogd -n -c3
startsecs = 2
stopwaitsecs = 5
